# proxy/snowflake.openbsd
#!/bin/ksh

daemon="/usr/local/bin/snowflake"
daemon_logger="daemon.info"
daemon_user="nobody"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
